# etc/systemd/system/noip.service
[Unit]
Description=Update dynamic DNS managed by no-ip.com

[Service]
ExecStart=/usr/sbin/noip

[Install]
WantedBy=multi-user.target
